# pqclid.service
[Unit]
Description=Progress Quest CLI (curses) in tmux
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/tmux-pq-supervisor

Restart=always
RestartSec=2

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
